# ami/salt/network/files/cpufreq-performance.service
[Unit]
Description=regnant CPU governor and IRQ pinning
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/regnant-irq-pin.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
